# rootfs/etc/services.d/2inputGrabber/run
#!/usr/bin/with-contenv bash

exec /bin/bash /app/bbuddy/example/grabInput.sh /dev/input/by-id/usb-0581_011c-event-kbd
